# deploy/oracle-cloud/pocketify-bot.service
[Unit]
Description=Pocketify Telegram signal alert bot
After=network.target

[Service]
Type=simple
User=pocketify
Group=pocketify
WorkingDirectory=/opt/pocketify
EnvironmentFile=/opt/pocketify/.env
ExecStart=/opt/pocketify/venv/bin/python signal_alert_bot.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
